# src/build_data/linux/installer/DEBIAN/postrm
#!/bin/sh

# https://wiki.debian.org/MaintainerScripts

set -e

if [ "$1" = "remove" ]; then
	if [ ! -d /run/systemd/system ]; then
		rm /etc/init.d/brickd
		update-rc.d brickd remove
	fi
fi

if [ "$1" = "purge" ]; then
	rm -f /var/log/brickd.log*
	rm -f /var/run/brickd.pid
fi

exit 0
