# src/build_data/alpine/brickd/brickd.initd
#!/sbin/openrc-run

command="/usr/sbin/brickd"
command_args="--daemon"
pidfile="/run/brickd.pid"

depend() {
	need localmount
	need net
	after bootmisc
}
