# src/build_data/apline/brickd.initd
#!/sbin/openrc-run

command="/usr/sbin/brickd"
command_args="--daemon"
pidfile="/var/run/brickd.pid"

depend() {
	need localmount
	need net
	after bootmisc
}
